# ignite/etc/sv/eth0-static/finish
#!/bin/sh
exec /sbin/ip link set dev eth0 down
